# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
autorestart=true

[program:r2t-editor]
command=/bin/bash -c "nodejs server.js"
directory=/opt/real-2-time-editor
autorestart=true